# onbuild/configs/supervisord.d/ambari-server.conf
[program:ambari-server]
command=ambari-server start
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
autostart=true
autorestart=false
priority=2
