# onbuild/configs/supervisord.d/deploy.conf
[program:AmbariDeploy]
command=/deploy/hdp-env.sh deploy_all
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stdout
stdout_logfile_maxbytes=0
autostart=true
autorestart=false
priority=50
